# roleprofile_modules/profile/templates/els-user-bootstrap-sync.cron.erb
# Invoke els-user-bootstrap to gather new DC1 users into projects
*/15 * * * * root /usr/bin/els-user-bootstrap --config /etc/els-openstack-integrations/els-openstack-integrations.conf
